# config/systemd/user/blueman-applet.service
[Unit]
Description=Bluetooth Manager
After=awesome.service

[Service]
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=graphical-session.target
