# config/systemd/user/conky.service
[Unit]
Description=Desktop Stats Display for X
#After=awesome.service
After=nitrogen.service
Requires=xfsettingsd.service

[Service]
#ExecStartPre=/bin/sleep 5
ExecStart=%h/Profile/bin/conky

[Install]
WantedBy=desktop.target
