# config/systemd/user/cryptkeeper.service
[Unit]
Description=ECryptFS Manager
#After=awesome.service

[Service]
ExecStart=/usr/bin/cryptkeeper

[Install]
WantedBy=desktop.target
