# config/systemd/user/default.target
[Unit]
Description=Default Target
Wants=graphical-session.target
AllowIsolate=true

[Install]
Alias=default.target
