# config/systemd/user/fdpowermon.service
[Unit]
Description=Power Manager Applet
After=awesome.service

[Service]
ExecStart=/usr/bin/fdpowermon

[Install]
WantedBy=graphical-session.target
