# config/systemd/user/hopsession.target
[Unit]
Description=Hop Session

[Install]
#Alias=default.target
WantedBy=graphical-session.target
