# config/systemd/user/keynav.service
[Unit]
Description=Shortcut Manager
WantedBy=awesome.service

[Service]
ExecStart=/usr/bin/keynav
