# config/systemd/user/mopidy.socket
[Socket]
ListenStream=6602

[Install]
WantedBy=hopsession.target
