# config/systemd/user/mpd.socket
[Socket]
ListenStream=6601

[Install]
WantedBy=hopsession.target
