# config/systemd/user/pasystray.service
[Unit]
Description=PulseAudio Tray Manager
Wants=awesome.service
After=awesome.service
Wants=pulseaudio.service
After=pulseaudio.service
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/pasystray

[Install]
WantedBy=graphical-session.target
