# config/systemd/user/revpn.service
[Unit]
Description=Revpn Script
After=awesome.service

[Service]
ExecStart=%h/Profile/bin/revpn.sh

[Install]
WantedBy=hopsession.target
