# config/systemd/user/safeeyes.service
[Unit]
Description=Safeeyes Time Manager
After=awesome.service

[Service]
ExecStart=/usr/bin/safeeyes

[Install]
WantedBy=graphical-session.target
