# config/systemd/user/workrave.service
[Unit]
Description=Time Manager
#After=awesome.service

[Service]
ExecStart=/usr/bin/workrave
