# config/systemd/user/xfsettingsd.service
[Unit]
Description=X Settings Manager for XFCE
PartOf=graphical-session.target

[Service]
Restart=on-failure
#RestartSec=5s
#Type=dbus
#BusName=org.xfce.Xfconf
#Environment=XAUTHORITY=%h/.Xauthority
ExecStart=/usr/bin/xfsettingsd --no-daemon --disable-wm-check --sm-client-disable

[Install]
WantedBy=desktop.target

