# config/systemd/user/xscreensaver.service
[Unit]
Description=XScreenSaver
#Before=awesome.service

[Service]
ExecStart=xscreensaver -nosplash

[Install]
WantedBy=desktop.target
