# config/systemd/user/xscreensaver.service
[Unit]
Description=XScreenSaver

[Service]
ExecStart=xscreensaver -nosplash

[Install]
WantedBy=graphical-session.target
